# services/albumplayer.service
[Unit]
Description=Album Player NFC Service
After=network-online.target bluetooth.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -u __INSTALL_DIR__/AlbumPlayer.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
