# services/webapp.service
[Unit]
Description=Album Player Web Interface
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dyonak
WorkingDirectory=/home/dyonak/album-player
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/dyonak/album-player/Webapp.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
