# services/webapp.service
[Unit]
Description=Album Player Web Interface
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 __INSTALL_DIR__/Webapp.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
