# services/wificonnect.service
[Unit]
Description=Album Player WiFi Provisioning
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -u __INSTALL_DIR__/wificonnect.py
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
