# services/wificonnect.service
[Unit]
Description=Album Player Wifi Config
Wants=network-online.target
After=network-online.target

[Service]
WorkingDirectory=/home/USERNAME/PifiConnector
# Optional environment file for fallback WiFi configuration overrides
# The '-' prefix means systemd won't fail if the file doesn't exist
EnvironmentFile=-/home/USERNAME/PifiConnector/wifi_fallback.env
ExecStart=/home/USERNAME/PifiConnector/wifivenv/bin/python3 wificonnect.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
